# files/etc/crontabs/root
* * * * * /root/connectionCheckAndConnect.sh > /dev/null
*/2 * * * * /root/sendStatus.sh > /dev/null
